# api-parser/crontab
*/1 * * * * python /usr/src/app/rest-parse-createWebsite.py -u #>> /usr/src/app/logfile